# web2py/applications/duty_calls/cron/crontab
#crontab

*/2  *  *  *  *  root *update/all
